# clippy.service
[Unit]
Description=clippy
After=network.target

[Service]
ExecStart=/usr/local/bin/clippy
Restart=on-failure

[Install]
WantedBy=default.target
